# resources/server_config/frontend/debian7_wheezy/public_site/etc/cron.d/videosquare-dev
# Cron jobs for: converter stream.videosquare.eu (dev)

MAILTO=hiba@videosqr.com
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

## Jobs

# Watcher (every 5 minutes)
*/5 * * * *		conv    /usr/bin/php -f /var/www/dev.videosquare.eu/modules/Jobs/watcher.php &> /dev/null

# Remove files (every 2 minutes)
* * * * *		conv    /usr/bin/php -f /var/www/dev.videosquare.eu/modules/Jobs/job_remove_files.php &> /dev/null

# Daily maintenance (every day 0:05)
5 0 * * *		conv    /usr/bin/php -f /var/www/dev.videosquare.eu/modules/Jobs/job_maintenance.php &> /dev/null

# Integrity check (Sunday 23:55)
55 23 * * 6		conv    /usr/bin/php -f /var/www/dev.videosquare.eu/modules/Jobs/job_integrity_check.php &> /dev/null
